# etc/profile.d/qemu-scripts.01.path.sh.in
case "::${PATH}::" in
*:"@package_installation_dir_unversioned@/bin":*)
	;;
*)
	PATH="${PATH}${PATH:+:}@package_installation_dir_unversioned@/bin"
	;;
esac
